# parking-predictor-24h.timer
[Unit]
Description=Run 24h parking predictor every hour (offset from scraper)

[Timer]
OnCalendar=*:01
AccuracySec=1s

[Install]
WantedBy=timers.target
